# com.tonyswu.ccckickstart.TASKNAME.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tonyswu.ccckickstart.TASKNAME</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/tonyswu/bin/ccckickstart.sh</string>
		<string>-name</string>
		<string>TASKNAME</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Day</key>
		<integer>1</integer>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>20</integer>
	</dict>
</dict>
</plist>
